# containers/docker/images/mesos/slave/mesos_executors.slice
[Unit]
Description=Mesos Executors Slice